# templates/webhook.service.erb
[Unit]
Description=R10K Webhook Service
After=syslog.target network.target

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/webhook
User=<%= @user %>
TimeoutStartSec=90
TimeoutStopSec=30

ExecStart=/usr/local/bin/webhook

KillMode=process

StandardOutput=syslog

[Install]
WantedBy=multi-user.target
