# templates/webhook.service.erb
[Unit]
Description=R10K Webhook Service
After=syslog.target network.target

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/webhook
RuntimeDirectory=webhook
User=<%= @user %>
PIDFile=/var/run/webhook/webhook.pid
TimeoutStartSec=90
TimeoutStopSec=30
RestartSec=10000

ExecStart=/usr/local/bin/webhook

KillMode=process

StandardOutput=syslog

[Install]
WantedBy=multi-user.target
