# templates/webhook.init.erb
#!/bin/bash
# webhook       Init script for running the r10k webhook daemon
#
# Author:       Zack Smith <zack.smith@puppetlabs.com>
#
# chkconfig: 2345 98 02
#
# description: Enables web based runs of r10k i.e. github enterprise post hooks
# processname: webhook

[ -f /etc/sysconfig/webhook ] && . /etc/sysconfig/webhook

# This is for ubuntu 14.04 (/var/run is mount to tmpfs)
if [ ! -d /var/run/webhook ]; then
  mkdir -p /var/run/webhook
  chown <%= @user %> /var/run/webhook
fi

pidfile='/var/run/webhook/webhook.pid'
webhook='/usr/local/bin/webhook'
DAEMON_USER='<%= @user %>'
RETVAL=0

# Source function library.
. <%= scope.lookupvar('r10k::params::functions_path') %>

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
export PATH

start() {
    echo
    echo -n $"Starting webhook: "
    <% if @osfamily == 'Debian' -%>
    start-stop-daemon --start --user ${DAEMON_USER:?} <%= scope.lookupvar('r10k::params::start_pidfile_args') %> --startas=$webhook
    <% else -%>
    daemon --user ${DAEMON_USER:?} <%= scope.lookupvar('r10k::params::start_pidfile_args') %> $webhook
    <% end -%>
    RETVAL=$?
    return $RETVAL
}

stop() {
    echo -n $"Stopping webhook: "
    pid=`netstat -lnp | awk  '/.*:8088.*ruby/{sub("/.*","",$NF);print $NF}'`
    #killproc -p $pidfile $webhook]
    if kill -9 $pid ; then
      success 'webhook'
      return 0
    else
      failure 'webhook'
      return 1
    fi
}

restart() {
    stop
    start
}

status() {
    if [ -f "${pidfile}" ]; then
        status -p $pidfile $webhook
        RETVAL=$?
    else
        RETVAL=3
    fi
    return $RETVAL
}


case "$1" in
    start)
        start
    ;;
    stop)
        stop
    ;;
    restart)
        restart
    ;;
    status)
        status
    ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart}"
        exit 1
esac

exit $RETVAL
