# templates/webhook.init.erb
#!/bin/bash
# webhook       Init script for running the r10k webhook daemon
#
# Author:       Zack Smith <zack.smith@puppetlabs.com>
#
# chkconfig: 2345 98 02
#
# description: Enables web based runs of r10k i.e. github enterprise post hooks
# processname: webhook
#
### BEGIN INIT INFO
# Required-Start:    $remote_fs $network
# Required-Stop:     $remote_fs $network
# Provides:          webhook
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Init script for running the r10k webhook daemon
### END INIT INFO

[ -f /etc/sysconfig/webhook ] && . /etc/sysconfig/webhook

# This is for ubuntu 14.04 (/var/run is mount to tmpfs)
if [ ! -d /var/run/webhook ]; then
  mkdir -p /var/run/webhook
  chown <%= @user %> /var/run/webhook
fi

pidfile='/var/run/webhook/webhook.pid'
webhook='/usr/local/bin/webhook'
DAEMON_USER='<%= @user %>'
RETVAL=0

# Source function library.
. <%= scope.lookupvar('r10k::params::functions_path') %>

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
export PATH

start() {
    echo
    echo -n $"Starting webhook: "
    <% if @osfamily == 'Debian' -%>
    start-stop-daemon --start --chuid ${DAEMON_USER} --user ${DAEMON_USER:?} <%= scope.lookupvar('r10k::params::start_pidfile_args') %> --startas=$webhook
    <% elsif @osfamily == 'SUSE' -%>
    start_daemon -u $DAEMON_USER -p $pidfile -f $webhook
    <% else -%>
    daemon --user ${DAEMON_USER:?} <%= scope.lookupvar('r10k::params::start_pidfile_args') %> $webhook
    <% end -%>
    RETVAL=$?
    return $RETVAL
}

stop() {
    echo -n $"Stopping webhook: "
    killproc -p $pidfile $webhook
    RETVAL=$?
    return $RETVAL
}

restart() {
    stop
    start
}

whstatus() {
    if [ -f "${pidfile}" ]; then
    <% if @osfamily == 'Debian' -%>
        status_of_proc -p $pidfile $webhook
    <% elsif @osfamily == 'SUSE' -%>
        checkproc -k -p $pidfile $webhook
    <% else -%>
        status -p $pidfile $webhook
    <% end -%>
        RETVAL=$?
    else
        RETVAL=3
    fi
    return $RETVAL
}


case "$1" in
    start)
        start
    ;;
    stop)
        stop
    ;;
    restart)
        restart
    ;;
    status)
        whstatus
    ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart}"
        exit 1
esac

exit $RETVAL
