# templates/webhook.service.erb
[Unit]
Description=R10K Webhook Service
After=syslog.target network.target

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/webhook
User=<%= @user %>
ExecStart=/usr/local/bin/webhook
PIDFile=/var/run/webhook/webhook.pid

[Install]
WantedBy=multi-user.target
